# conf/etc/service/graphite/run
#!/bin/bash

exec /usr/local/bin/gunicorn --chdir /opt/graphite/webapp/graphite graphite_wsgi:application --bind 127.0.0.1:8080
